# trunk/contrib/osx/net.freehaven.tor.plist.in
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>net.freehaven.tor</string>

	<key>Program</key>
	<string>/Library/Tor/tor</string>

	<key>ProgramArguments</key>
	<array>
	<string>tor</string>
	</array>

	<key>UserName</key>
	<string>_tor</string>

	<key>GroupName</key>
	<string>daemon</string>

	<key>OnDemand</key>
	<false/>

	<key>RunAtLoad</key>
	<true/>

	<key>ServiceDescription</key>
	<string>Tor @VERSION@  http://tor.eff.org</string>

	<key>LowPriorityIO</key>
	<true/>

	<key>HardResourceLimits</key>
	<dict>
	<key>NumberOfFiles</key>
	<integer>8192</integer>
	</dict>
</dict>
</plist>
